# noaa-role/templates/noaa.service.j2
[Unit]
Description=Stream NOAA Weather Radio to Icecast

[Service]
User={{ ansible_ssh_user }}
Type=simple
WorkingDirectory=/opt/noaa
ExecStart=/usr/bin/liquidsoap /opt/noaa/noaa-stream.liq
